# hometime.service
[Unit]
Description=Hometime Calendar LED Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/ndomino/hometime
ExecStart=node dist/index.js
Restart=always
RestartSec=10
Environment=GCALCLI_CONFIG=/opt/ndomino/hometime/config.toml
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
SyslogIdentifier=hometime

[Install]
WantedBy=multi-user.target
